# lima-init.openrc
#!/sbin/openrc-run

depend() {
  after lima-init-local
  after net
  after procfs
  after qemu-binfmt
  provide lima-init
}

start() {
  lima-init
  eend 0
}
