# lima-init-local.openrc
#!/sbin/openrc-run

depend() {
  after localmount
  before net
  provide lima-init-local
}

start() {
  lima-init --local
  eend 0
}
