# lima-init.openrc
#!/sbin/openrc-run

depend() {
  after lima-init-local
  after net
  provide lima-init
}

start() {
  lima-init
  eend 0
}
